# hw2/activate.service
[Unit]
Description=activate led
RemainAfterExit=true
[Service]
Type=oneshot
ExecStart=/usr/bin/activate_led
