# hw2/deactivate.service
[Unit]
Description=Deactivate service

[Service]
Type= oneshot

ExecStart=/usr/bin/deactivate_led